# repo/pkg/deb/skarphed-repo/DEBIAN/postrm
#!/bin/sh

set -e

update-rc.d skdrepo remove
